# target/linux/layerscape/base-files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
if grep -qs "LS1012A" /sys/devices/soc0/family
then
modprobe pfe.ko
fi
exit 0
